# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:rdcd]
command=/opt/rocm-5.5.0/bin/rdcd -u

[program:prometheus]
command=python3 /rdc/python_binding/rdc_prometheus.py
